# systemd_services/failsafe.service
[Unit]
Description=DataClaw Emergency Failsafe Trigger
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/dataclaw/scripts/trigger_safe_mode.sh
